# debian/ebeam-diva.prerm
#!/bin/sh

killall eBeam_Stylus_Preferences

invoke-rc.d ebeam-stylus-driver stop

exit 0
