# debian/postrm
#!/bin/sh

update-rc.d ebeam-stylus-driver remove

exit 0

