# debian/postinst
#!/bin/sh


update-rc.d ebeam-stylus-driver defaults 15

invoke-rc.d ebeam-stylus-driver start

exit 0
