# debian/prerm
#!/bin/sh

invoke-rc.d ebeam-stylus-driver stop

exit 0
