# etc/xdg/autostart/ebeam-stylus-preferences.desktop
[Desktop Entry]
Name=eBeam Stylus Preferences
Exec=/opt/diva/delayed-preferences.sh
Type=Application
Terminal=false
Categories=Other;
